# core/ops/snort-fetch.timer
[Unit]
Description=Run Snort fetch periodically

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
Unit=snort-fetch.service

[Install]
WantedBy=timers.target

